# crontab
# m h  dom mon dow   command
05 12 * * * TP1_BACKUP.sh 
